# image-builder/overlay/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin meteohub --noclear %I $TERM
Type=idle