# image-builder/overlay/etc/systemd/system/kegweather.service
[Unit]
Description=KEGWeatherHist

[Service]
ExecStart=/home/alarm/start

[Install]
WantedBy=multi-user.target